# build/crontab
# SCAN ALL IMAGES EVERY SUNDAY AT MIDNIGHT
0 0 1 * * /scripts/scan-all.sh >/dev/null 2>&1
